# metody/personal-ai-assistant/implementation/stage-0.6/personal-assistant-bot.service
[Unit]
Description=Personal AI Assistant — Stage 4 Bot (callbacks + voice receiver)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=personal-assistant
Group=personal-assistant
WorkingDirectory=/opt/personal-assistant

EnvironmentFile=/etc/personal-assistant/config.env

LoadCredentialEncrypted=bot-token:/etc/credstore.encrypted/bot-token.cred

ExecStart=/usr/bin/node /opt/personal-assistant/bot.mjs

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-assistant-bot

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/personal-assistant /var/log/personal-assistant
PrivateTmp=true
MemoryMax=300M

[Install]
WantedBy=multi-user.target
