# metody/personal-ai-assistant/implementation/stage-3/personal-assistant-draft-gen.service
[Unit]
Description=Personal AI Assistant — Stage 3 Draft Generator Worker
After=network-online.target personal-assistant-listener.service personal-assistant-triage.service
Wants=network-online.target
Requires=personal-assistant-listener.service personal-assistant-triage.service

[Service]
Type=simple
User=personal-assistant
Group=personal-assistant
WorkingDirectory=/opt/personal-assistant
EnvironmentFile=/opt/personal-assistant/.env
ExecStart=/usr/bin/node /opt/personal-assistant/draft_gen.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pa-draft-gen

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
ReadWritePaths=/opt/personal-assistant
MemoryMax=400M
TasksMax=50

[Install]
WantedBy=multi-user.target
