# metody/personal-ai-assistant/implementation/stage-2/personal-assistant-triage.service
[Unit]
Description=Personal AI Assistant — Triage classifier
After=network-online.target personal-assistant-listener.service
Wants=network-online.target

[Service]
Type=simple
User=personal-assistant
Group=personal-assistant
WorkingDirectory=/opt/personal-assistant
EnvironmentFile=/opt/personal-assistant/.env
ExecStart=/usr/bin/node /opt/personal-assistant/triage.mjs
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-assistant-triage

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/personal-assistant /var/log/personal-assistant
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryMax=300M
TasksMax=64

[Install]
WantedBy=multi-user.target
